# cron.d/obs-ci-integration
# mingw32-kmymoney5
0 0 * * * /usr/local/bin/osc-git-trigger
